# dist/siance-backend.service
[Unit]
Description=Siance Backend

[Service]
Type=simple
ExecStart=/bin/bash "PATH OF siance-backend-start.sh"

[Install]
WantedBy=multi-user.target
